# .local/share/noteworthy/evremap/etc/init.d/evremap
#!/sbin/openrc-run

supervisor=supervise-daemon
command="/usr/local/bin/evremap"
command_args="remap /etc/conf.d/evremap.toml"

